# host-config/linux/etc/cron.d/backup
SHELL=/bin/bash
15 2 * * * root /root/backup.sh
